# supervisor.conf
[supervisord]
nodaemon = true
[program:gunicorn]
command=/usr/local/bin/gunicorn --chdir /usr/src/app --workers 3 --bind unix:/usr/src/app/app.sock moove_doc_signing.wsgi:application
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
[program:nginx]
command = /usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0